# packaging/macos/dev.sakimori.watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.sakimori.watch</string>

  <!-- EDIT: absolute path to the installed `sakimori` binary and the
       directory/directories you want watched. Multiple roots can be
       added as additional <string> entries in ProgramArguments. -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/sakimori</string>
    <string>deps</string>
    <string>watch</string>
    <string>--min-age</string>
    <string>7d</string>
    <string>--notifier</string>
    <string>mac</string>
    <string>/Users/YOU/code</string>
  </array>

  <!-- Start on login and restart automatically if it ever crashes. -->
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <!-- Log to ~/Library/Logs so Console.app picks them up. -->
  <key>StandardOutPath</key>
  <string>/tmp/sakimori-watch.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sakimori-watch.log</string>

  <!-- So osascript can post notifications. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>RUST_LOG</key>
    <string>info</string>
  </dict>
</dict>
</plist>
